# roles/autorestic/files/autorestic.service
[Unit]
Description=Autorestic cron

[Service]
Environment=TZ=UTC
Type=oneshot
ExecStart=/etc/autorestic/autorestic

[Install]
WantedBy=timers.target
