# roles/openclaw/files/claude-max-api-proxy.service
[Unit]
Description=Claude Max API Proxy
After=network.target

[Service]
Type=simple
User=vik
EnvironmentFile=/mnt/tank/openclaw/.openclaw/workspace/.env
ExecStart=/usr/bin/claude-max-api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
