# roles/openclaw/files/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network.target claude-max-api-proxy.service
Requires=claude-max-api-proxy.service

[Service]
Type=simple
User=vik
Environment=OPENCLAW_HOME=/mnt/tank/openclaw
ExecStart=/usr/bin/openclaw gateway
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
