# roles/zfs/files/syncoid.timer
[Unit]
Description=Syncoid pull every night
Requires=syncoid.service

[Timer]
OnCalendar=05:05:05
Persistent=true

[Install]
WantedBy=timers.target
